# profile.d/bat.sh
# shellcheck shell=bash

# bat config
# set light theme
export BAT_THEME="OneHalfLight"
