# profile.d/bat.sh
# bat config
# set light theme
export BAT_THEME="OneHalfLight"
