# profile.d/editor.sh
# shellcheck shell=bash

export EDITOR=nano
