# profile.d/osx.sh
if [ -e "$HOME"/dotfiles/ensurepath.sh ]; then
    source "$HOME"/dotfiles/ensurepath.sh
    ensurepath /usr/local/bin
fi
